# systemd/the-bot-celery-beat.service
[Unit]
Description=THE BOT Platform - Celery Beat Scheduler Service
Documentation=https://docs.celeryproject.io/
After=network.target redis.service celery-worker.service
Wants=redis.service celery-worker.service

[Service]
Type=simple
User=www-data
Group=www-data

# Working directory
WorkingDirectory=/home/mg/THE_BOT_platform/backend

# Activate virtual environment and start Celery beat
ExecStart=/home/mg/THE_BOT_platform/.venv/bin/celery -A core \
    beat \
    --loglevel=info \
    --logfile=/var/log/celery/beat.log \
    --pidfile=/var/run/celery/beat.pid \
    --scheduler=django_celery_beat.schedulers:DatabaseScheduler

# Graceful restart
ExecReload=/bin/kill -s HUP $MAINPID

# Environment variables
Environment="DJANGO_SETTINGS_MODULE=config.settings"
Environment="PYTHONUNBUFFERED=1"
EnvironmentFile=/home/mg/THE_BOT_platform/.env

# Restart policy
Restart=always
RestartSec=10

# Process management
KillMode=process
KillSignal=SIGTERM

# Security
NoNewPrivileges=yes
PrivateTmp=yes

# Timeouts
TimeoutStartSec=300
TimeoutStopSec=60

[Install]
WantedBy=multi-user.target
