# scripts/deployment/thebot-celery-beat.service
[Unit]
Description=THE_BOT Platform - Celery Beat Scheduler
After=network.target redis.service
Wants=redis.service

[Service]
Type=simple
User=mg
Group=mg
WorkingDirectory=/home/mg/THE_BOT_platform/backend
Environment="PATH=/home/mg/venv/bin"
Environment="DJANGO_SETTINGS_MODULE=config.settings"
Environment="PYTHONUNBUFFERED=1"

# Create logs directory
ExecStartPre=/bin/mkdir -p /home/mg/THE_BOT_platform/backend/logs

# Start command
ExecStart=/home/mg/venv/bin/celery -A config.celery \
    beat \
    --loglevel=info \
    --logfile=/home/mg/THE_BOT_platform/backend/logs/celery_beat.log \
    --pidfile=/home/mg/THE_BOT_platform/backend/logs/celery_beat.pid \
    --scheduler django_celery_beat.schedulers:DatabaseScheduler

# Restart policy
Restart=on-failure
RestartSec=10
StartLimitInterval=600
StartLimitBurst=3

# Process management
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=15

# Security
PrivateTmp=yes
NoNewPrivileges=yes

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=thebot-celery-beat

[Install]
WantedBy=multi-user.target
